# infra/single-vm/systemd/a11y-dashboard.service
[Unit]
Description=A11y Audit Tool — dashboard (Next.js)
After=network.target a11y-api.service
Wants=a11y-api.service

[Service]
Type=simple
User=a11y
Group=a11y
WorkingDirectory=/opt/a11y/repo/apps/dashboard
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
Environment=API_BASE=http://127.0.0.1:3001

ExecStart=/usr/bin/env pnpm --filter @a11y/dashboard start
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/a11y/dashboard.log
StandardError=append:/var/log/a11y/dashboard.log

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/a11y /opt/a11y/repo/apps/dashboard/.next
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
